# deployment/roles/wildlex/templates/wildlex_prod.cron.j2
# WILDLEX production cron
* * * * * php curl -s http://www.wildlex.org/cron/{{ wildlex_cron_key }}
